# oga-events/oga-events.service
[Unit]
Description=odroidgo2 global events listener

[Service]
ExecStart=/usr/bin/oga-events
Type=simple

[Install]
WantedBy=multi-user.target
