# ogage/ogage.service
[Unit]
Description=ogage

[Service]
ExecStart=/usr/bin/ogage
Type=simple

[Install]
WantedBy=multi-user.target
