# system-odroidgo2/onboot.service
[Unit]
Description=Onboot

[Service]
Type=forking
ExecStart=/usr/local/sbin/onboot.sh

[Install]
WantedBy=multi-user.target
